# ROOT/home/aidan/.config/systemd/user/daily-save-docs.timer
[Unit]
Description=Run Daily Git Commit at the End of Day

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
