# ROOT/home/aidan/.config/systemd/user/save-docs.service
[Unit]
Description=Save Documents folder by commiting and pushing to remote repository.

[Service]
Type=oneshot
ExecStart=/home/aidan/bin/save-docs
