# HOME/.config/systemd/user/daily-save-docs.service
[Unit]
Description=Daily Git Commit for Documents Folder

[Service]
Type=oneshot
ExecStart=/home/aidan/bin/save-docs
